# src/Bot_Base.service
[Unit]
Description=Discord Bot for unionrolistes.fr
RequiresMountsFor=/usr/local/bin/
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/URbot/start.py
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=inherit

[Install]
WantedBy=multi-user.target